# packages/digital-hydrant/digital_hydrant-0.1.1-py3-none-any.whl/digital_hydrant/systemd/digital-hydrant.service
[Unit]
Description=Digital Hydrant

[Service]
Type=forking
ExecStart=/usr/bin/hydrant

[Install]
WantedBy=multi-user.target
